# local_modules/supervisor/templates/supervisord.conf.erb
[inet_http_server]
port = 0.0.0.0:<%= @supervisor_port %>

[supervisord]
logfile = <%= @dir %>/supervisord.log
logfile_maxbytes = 100MB
logfile_backups = 10
loglevel = info
pidfile = <%= @dir %>/supervisord.pid
nodaemon = false

[supervisorctl]
serverurl = http://localhost:<%= @supervisor_port %>
username = <%= @user %>
password = b

[program:readthedocs]
directory = <%= @checkouts %>/readthedocs
command = <%= @venv %>/bin/python manage.py runserver 0.0.0.0:<%= @port %>
autostart = true
autorestart = true
startretries = 3
redirect_stderr = true								; rtd logs to stderr - weird
stderr_logfile = /home/<%= @user %>/rtd.err.log
stdout_logfile = /home/<%= @user %>/rtd.out.log		; would be empty if not redirect_stderr
user = <%= @user %>
stopasgroup = true  								; so child process is stopped